# case/daily_close_brief.conf
[program:daily_close_brief]
# 优化：使用 Python 解释器绝对路径，避免环境版本问题
command=/root/.pyenv/shims/python /home/fantom/datastore/case/daily_close_brief.py --scheduled
# 工作目录，保持不变
directory=/home/fantom/datastore/case
# 修正：直接指定具体用户，消除环境变量依赖
user=root
# 自动启动/重启，保持不变
autostart=true
autorestart=true
# 启动相关配置，保持不变
startsecs=10
startretries=3
# 停止相关配置，保持不变
stopwaitsecs=10
stopsignal=INT
stopasgroup=true
killasgroup=true
# 错误输出重定向，保持不变
redirect_stderr=true
# 日志配置，保持不变
stdout_logfile=/home/fantom/datastore/case/logs/supervisor/daily_close_brief.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10