# case/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:daily_close_brief]
# 优化：使用 Python 解释器绝对路径，避免环境版本问题
command=/usr/bin/python3 daily_close_brief.py --scheduled
# 工作目录，保持不变
directory=/home/fantom/datastore/case
# 修正：直接指定具体用户，消除环境变量依赖
user=fantom
# 自动启动/重启，保持不变
autostart=true
autorestart=true
# 启动相关配置，保持不变
startsecs=10
startretries=3
# 停止相关配置，保持不变
stopwaitsecs=10
stopsignal=INT
stopasgroup=true
killasgroup=true
# 错误输出重定向，保持不变
redirect_stderr=true
# 日志配置，保持不变
stdout_logfile=/home/fantom/datastore/case/logs/supervisor/daily_close_brief.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
# 优化：环境变量配置（格式规范，后续可扩展多个变量）
environment=PYTHONPATH="/home/fantom/datastore/case"