# deploy/celery-worker.service
[Unit]
Description=Celery Worker for Qlib Training
After=network.target redis-server.service

[Service]
Type=simple
User=fantom
WorkingDirectory=/home/fantom/datastore/apps/api
ExecStart=/home/fantom/.pyenv/versions/3.12.0/bin/celery -A app.celery_app worker --loglevel=info --concurrency=1 --pool=solo
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
